# claude-code/launchd/com.claude.sleep-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.sleep-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME/ghq/github.com/DaichiHoshina/ai-tools" &amp;&amp; /Users/daichi.hoshina/.local/bin/claude -p "夜間 pipeline を実行する。sleep-review skill の staging 手順に従い、改善提案を stage 保存する。翌朝 /sleep-review triage で読める形にする。"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/daichi.hoshina/.local/bin</string>
        <key>HOME</key>
        <string>/Users/daichi.hoshina</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>3</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/daichi.hoshina/.claude/logs/launchd/sleep-review.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/daichi.hoshina/.claude/logs/launchd/sleep-review.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
